# package/base-files/files/etc/uci-defaults/99_disable-qos-by-default
#!/bin/sh

/etc/init.d/qos_gargoyle disable

exit 0