# package/pifii/ac-discovery-server/src/ac-discovery-server.init
#!/bin/sh /etc/rc.common

START=94

start(){
	service_start /usr/sbin/ac-discovery-server
}

stop(){
	service_stop /usr/sbin/ac-discovery-server	
}

restart(){
        stop
        start
}
