# package/base-files/files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
test -f /sys/module/nf_conntrack/parameters/hashsize && echo 64000 > /sys/module/nf_conntrack/parameters/hashsize
exit 0
